# nspawn/system/systemd-nspawn@.service.d/override.conf
[Service]
ExecStartPost=/usr/local/etc/nspawn/machine.sh %I
